# deploy/systemd/polypulse-monitor.service
[Unit]
Description=PolyPulse persistent market monitor
Documentation=/home/PolyPulse/deploy/README.md
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/PolyPulse
EnvironmentFile=/home/PolyPulse/.env
Environment=POLYPULSE_HOME=/home/PolyPulse
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=256
ExecStartPre=/home/PolyPulse/deploy/scripts/healthcheck.sh --preflight
ExecStart=/usr/bin/env bash -lc 'cd /home/PolyPulse; MODE="${POLYPULSE_EXECUTION_MODE:-paper}"; if [ "$MODE" = "live" ]; then if [ "${POLYPULSE_LIVE_CONFIRM:-}" != "LIVE" ]; then echo "live monitor requires POLYPULSE_LIVE_CONFIRM=LIVE"; exit 9; fi; exec /usr/bin/env node ./bin/polypulse.js monitor run --mode live --env-file /home/PolyPulse/.env --confirm LIVE --loop; fi; exec /usr/bin/env node ./bin/polypulse.js monitor run --mode paper --env-file /home/PolyPulse/.env --loop'
Restart=always
RestartSec=15
TimeoutStartSec=180
TimeoutStopSec=30
KillSignal=SIGTERM
NoNewPrivileges=true
MemoryMax=512M
CPUQuota=80%
StandardOutput=append:/home/PolyPulse/logs/polypulse-monitor.log
StandardError=append:/home/PolyPulse/logs/polypulse-monitor.err.log

[Install]
WantedBy=multi-user.target
